# modules/dc_rails/templates/clockwork.upstart.erb
start on runlevel [2345]
stop on runlevel [!2345]

setuid rails
setgid rails

respawn
respawn limit 3 30

normal exit 0 TERM
stop on (stopping workers or runlevel [06])

script
# this script runs in /bin/sh by default
# respawn as bash so we can source in rbenv
exec /bin/bash <<'EOT'
  source /etc/default/unicorn_stronghold
  cd /home/rails/stronghold/current
  exec /home/rails/.rbenv/shims/bundle exec clockwork clock.rb
EOT
end script