# modules/dc_rails/templates/sidekiq.upstart.erb
start on runlevel [2345]
stop on runlevel [!2345]

setuid rails
setgid rails

chdir <%= @app_home %>

env DB_PASSWORD=<%= @db_password %>
env RAILS_ENV=<%= @rails_env %>

exec <%= @home %>/.rbenv/shims/bundle exec sidekiq
